# gcp-worker/quantum-alchemy-gcp.service
[Unit]
Description=Quantum Alchemy Engine — GCP Worker (DFT + GNN)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/quantum-alchemy
EnvironmentFile=/etc/quantum-alchemy.env

# Pull latest code from git before each start so restarts always run HEAD
ExecStartPre=/usr/bin/git -C /opt/quantum-alchemy pull --ff-only

# tsx runs TypeScript directly — no build step needed
ExecStart=/usr/bin/npx tsx gcp-worker/index.ts

# Graceful shutdown: SIGTERM then SIGKILL after 35s
TimeoutStopSec=35
KillMode=control-group
KillSignal=SIGTERM

Restart=always
RestartSec=10

# Allow large GNN weight blobs in the journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-alchemy-gcp

# Increase open file limit for Postgres connections + QE tmp files
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
